# private_apps/boot-logger/boot-logger-shutdown.service
[Unit]
Description=Boot Logger Shutdown Marker
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/home/tha/apps/boot-logger/shutdown-marker.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=halt.target reboot.target shutdown.target
